# scripts/deploy-v2/combat-v2-direct.service
[Unit]
Description=Combat-v2 Backend + Frontend (:3001)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/combat-v2/apps/backend
Environment=PATH=/root/.nvm/versions/node/v22.22.3/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=COMBAT_API=http://localhost:3001
ExecStart=/root/.nvm/versions/node/v22.22.3/bin/npx tsx src/server.ts
Restart=always
RestartSec=5
StandardOutput=append:/opt/combat-v2/backend.log
StandardError=append:/opt/combat-v2/backend.log

[Install]
WantedBy=multi-user.target